# launchd/com.hermes.deeptutor-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.deeptutor-frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/can/DeepCode/DeepTutor/web && node --max-old-space-size=4096 ./node_modules/next/dist/bin/next dev -p 3782</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/can/DeepCode/DeepTutor/web/node_modules/.bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_OPTIONS</key>
        <string>--max-old-space-size=4096</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/can/.hermes/logs/deeptutor-frontend.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/can/.hermes/logs/deeptutor-frontend.error.log</string>
</dict>
</plist>
