# launchd/com.hermes.deeptutor-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.deeptutor-backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/can/.openharness-venv/bin/deeptutor</string>
        <string>serve</string>
        <string>--port</string>
        <string>8001</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/can/DeepCode/DeepTutor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/can/.openharness-venv/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/can/.hermes/logs/deeptutor-backend.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/can/.hermes/logs/deeptutor-backend.error.log</string>
</dict>
</plist>
